# deployment/systemd/narrowcasting-kiosk.service
[Unit]
Description=Narrowcasting Chromium Kiosk
After=graphical.target narrowcasting-player.service
Wants=graphical.target narrowcasting-player.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/narrowcasting
Environment=NARROWCASTING_ROOT=/opt/narrowcasting
Environment=KIOSK_URL=http://localhost:4174/player
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/narrowcasting/scripts/start-kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
